# init.d/updateHostname
#!/bin/sh /etc/rc.common

START=79

BIN=/usr/bin/updateHostname

boot() {
    $BIN 
}